# deployment/strava-backend.service
[Unit]
Description=Strava Activity Copilot Backend
After=network.target

[Service]
User=andy
WorkingDirectory=/var/www/strava-activity-copilot
EnvironmentFile=/var/www/strava-activity-copilot/backend/.env
ExecStart=/var/www/strava-activity-copilot/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
